# systemd/ubuntu-metrics-exporter.service
[Unit]
Description=Ubuntu Metrics Exporter (Prometheus)
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/ubuntu-metrics
ExecStart=/opt/ubuntu-metrics/venv/bin/python /opt/ubuntu-metrics/ubuntu_metrics_exporter.py
Restart=always
RestartSec=10
Environment="PYTHONUNBUFFERED=1"

# Log location (journalctl still works)
StandardOutput=append:/var/log/ubuntu_metrics_exporter.log
StandardError=append:/var/log/ubuntu_metrics_exporter.err

[Install]
WantedBy=multi-user.target
